# meta-common/recipes-phosphor/session-management/session-management/xyz.openbmc_project.SessionManager.service
[Unit]
Description=Session Management Service
Requires=bmcweb.service
After=bmcweb.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/session-management

[Install]
WantedBy=multi-user.target
